# rc.local
#! /bin/bash
/etc/scripts/static-arp.sh
systemctl restart firewall
exit 0
