# rc.local
#! /bin/bash
/etc/scripts/static-arp.sh
exit 0
